# etc/profile.d/opt.sh
# include /opt/bin into system paths

export PATH=$PATH:/opt/bin
